# daemon/deployments/dicompute-daemon.service
[Unit]
Description=DICOMPUTE Provider Daemon
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=dicompute
Group=dicompute
WorkingDirectory=/opt/dicompute/daemon
ExecStart=/opt/dicompute/daemon/daemon
Restart=always
RestartSec=5
Environment=LOG_LEVEL=info

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/dicompute/data

[Install]
WantedBy=multi-user.target
